# home/.config/systemd/user/devilspie2.service
[Unit]
Description=Devilspie2
After=gvfs-daemon.service

[Service]
ExecStart=/usr/bin/devilspie2

[Install]
WantedBy=default.target
